# contrib/com.github.zchee.imectl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.zchee.imectl</string>

    <!-- Adjust the path to where you installed imectl. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/imectl</string>
        <string>daemon</string>
    </array>

    <!-- Start at login and keep the warm-connection daemon alive. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/imectl-daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/imectl-daemon.err.log</string>
</dict>
</plist>
